# agent/installer/launch_agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.renderstudio.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/RenderStudio/renderstudio-agent</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/renderstudio-agent.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/renderstudio-agent-error.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>RS_API_BASE</key>
    <string>http://localhost:8000</string>
  </dict>
</dict>
</plist>
